# platform/broadcom/sonic-platform-modules-delta/debian/platform-modules-ag9032v1.postinst
#!/bin/sh
systemctl enable ag9032v1_platform_fancontrol.service
systemctl enable ag9032v1_platform_led.service
systemctl enable ag9032v1_platform_modules.service
